# vnc-hoster.service
[Unit]
Description=VNC Hoster

[Service]
Type=simple
WorkingDirectory=/vnc-hoster
ExecStart=python3 client.py

[Install]
WantedBy=multi-user.target
